# supervisord.conf
[supervisord]
nodaemon=true

[program:consulagentclient]
command=/bin/consulagentclient
