# supervisord.conf
[supervisord]
nodaemon=true

[program:consulagentclient]
command=/bin/consulagentclient

[include]
files= /etc/supervisor/conf.d/*.conf
