# supervisord.conf
[supervisord]
nodaemon=true

[program:consul]
command=/bin/start
